# config/crontab
0 */6 * * * /usr/local/bin/php /partkeepr/app/console partkeepr:cron:run >> /partkeepr/app/logs/cron.log
# Leave empty line at end of file
